# batterytator.service
[Unit]
Description=Battery monitoring service that send notifications

[Service]
ExecStart=/usr/bin/batterytator
